# deploy/telegram-mt5-bot.service
[Unit]
Description=Telegram MT5 Trading Bot
Documentation=https://github.com/your-repo/telegram-mt5-bot
After=network-online.target mt5-rpyc-server.service
Wants=network-online.target mt5-rpyc-server.service

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/telegram-mt5-bot
ExecStart=/opt/telegram-mt5-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment
Environment=DISPLAY=:0
EnvironmentFile=/opt/telegram-mt5-bot/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/telegram-mt5-bot
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-mt5-bot

[Install]
WantedBy=multi-user.target
